# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/busybox-mdev.postinst
	[ -n "$D" ] && OPT="-r $D" || OPT="-s"
	update-rc.d $OPT mdev start 06 S .
