# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/busybox-syslog.postrm
	if test "$1" = "remove" -o "$1" = "purge"; then
		if ! test -e "/etc/init.d/syslog"; then
			[ -n "$D" ] && OPT="-r $D" || OPT="-s"
			update-rc.d $OPT syslog remove
		fi
	fi
