# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/busybox.list
/etc/udhcpc.d/
/etc/busybox.links
/bin/sh
/share/udhcpc/default.script
/share/
/etc/
/etc/udhcpc.d/50default
/bin/busybox
/etc/default/hwclock
/etc/init.d/
/etc/init.d/busybox-cron
/bin/
/share/udhcpc/
/etc/init.d/hwclock.sh
/etc/default/
