# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/busybox-mdev.list
/etc/mdev/usb.sh
/etc/mdev/autoumountsdcard.sh
/etc/mdev/
/etc/mdev/automountsdcard.sh
/etc/mdev/find-touchscreen.sh
/etc/mdev.conf
/etc/init.d/mdev
