# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/busybox-syslog.postinst
	update-alternatives --install /etc/init.d/syslog syslog-init syslog.busybox 50

	[ -n "$D" ] && OPT="-r $D" || OPT="-s"
	# remove all rc.d-links potentially created from alternative
	# syslog packages before creating new ones
	update-rc.d $OPT -f syslog remove
	update-rc.d $OPT syslog start 20 2 3 4 5 . stop 90 0 1 6 .
