# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/data-ipa-cfg-mgr.postinst
[ -n "$D" ] && OPT="-r $D" || OPT="-s"
        update-rc.d $OPT -f start_ipacm_le remove
        update-rc.d $OPT start_ipacm_le start 90 2 3 4 5 . stop 10 0 1 6 .

