# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/data-ipa-cfg-mgr.list
/usr/lib/libipanat.so.1
/etc/init.d/start_ipacm_le
/usr/bin/ipacm
/etc/IPACM_cfg.xml
