# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/diag.list
/usr/bin/test_diag
/usr/bin/diag_uart_log
/etc/init.d/chgrp-diag
/usr/bin/PktRspTest
/usr/lib/libdiag.so.1
/usr/bin/diag_mdlog
/usr/bin/diag_klog
