# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.d/
/etc/dnsmasq.conf
/etc/default_config/dhcp_var
/usr/bin/dnsmasq
/etc/config/dhcp_var
/etc/init.d/dnsmasq
/usr/bin/dhcp_script.sh
