# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/ebtables.list
/usr/lib/libebt_802_3.so
/usr/lib/libebt_stp.so
/etc/rc.d/init.d/ebtables
/usr/lib/libebt_standard.so
/usr/lib/libebt_redirect.so
/usr/lib/libebt_nflog.so
/usr/lib/libebt_ulog.so
/usr/lib/libebt_vlan.so
/usr/sbin/ebtables
/usr/lib/libebtable_nat.so
/usr/lib/libebt_nat.so
/usr/lib/libebtc.so
/usr/lib/libebtable_filter.so
/usr/sbin/ebtables-restore
/usr/lib/libebt_arpreply.so
/usr/lib/libebt_log.so
/usr/lib/libebt_ip6.so
/etc/ebtables/ebtables-config
/etc/ebtables/
/usr/lib/libebtable_broute.so
/usr/lib/libebt_mark.so
/usr/lib/libebt_pkttype.so
/etc/ethertypes
/usr/lib/libebt_limit.so
/etc/rc.d/
/usr/lib/libebt_mark_m.so
/usr/sbin/ebtables-save
/usr/lib/libebt_among.so
/usr/lib/libebt_ip.so
/usr/lib/libebt_arp.so
/etc/rc.d/init.d/
