# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/lighttpd.list
/etc/lighttpd.user
/www/pages/dav/
/etc/init.d/lighttpd
/www/logs/
/etc/openssl.cnf
/etc/lighttpd.conf
/usr/sbin/lighttpd
/usr/sbin/lighttpd-angel
/www/var/
/www/pages/index.html
/www/pages/
