# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/ntpclient.list
/etc/zoneinfo/GMT60
/etc/zoneinfo/GMT240
/etc/zoneinfo/GMT420
/etc/zoneinfo/GMT180
/etc/localtime
/etc/zoneinfo/GMT-120
/etc/zoneinfo/GMT270
/etc/zoneinfo/GMT-210
/etc/zoneinfo/GMT600
/etc/zoneinfo/GMT360
/etc/zoneinfo/GMT-300
/etc/zoneinfo/
/etc/zoneinfo/GMT540
/etc/zoneinfo/GMT-240
/etc/zoneinfo/GMT-330
/etc/zoneinfo/GMT720
/etc/zoneinfo/GMT480
/etc/zoneinfo/GMT-420
/etc/zoneinfo/GMT-180
/etc/zoneinfo/GMT-345
/etc/zoneinfo/GMT0
/etc/zoneinfo/GMT-270
/etc/zoneinfo/GMT660
/etc/zoneinfo/GMT-600
/etc/zoneinfo/GMT-360
/etc/zoneinfo/GMT-540
/etc/zoneinfo/GMT-390
/etc/init.d/ntp_client
/etc/zoneinfo/GMT-720
/etc/zoneinfo/GMT-480
/etc/zoneinfo/GMT-570
/etc/zoneinfo/GMT-660
/usr/bin/ntp_client
/etc/zoneinfo/GMT-780
/etc/zoneinfo/GMT-60
/etc/zoneinfo/GMT120
/etc/zoneinfo/GMT210
/etc/zoneinfo/GMT300
