# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/powerapp.list
/etc/init.d/reset_reboot_cookie
/sbin/sys_shutdown
/etc/reboot-cookie
/sbin/reboot-recovery
/sbin/reboot-bootloader
/sbin/sys_reboot
/sbin/powerapp
