# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/qmi.list
/usr/lib/libqmiidl.so.1
/usr/lib/libqmi.so.1
/usr/bin/qmi_simple_ril_test
/etc/init.d/qmuxd
/usr/lib/libqmiservices.so.1
/usr/lib/libqcci_legacy.so.1
/usr/lib/libqmi_client_qmux.so.1
/usr/bin/qmuxd
