# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/ril.list
/etc/config/rild
/etc/default_config/rild
/usr/lib/libril.so
/usr/bin/rild
/etc/init.d/rild
