# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/run-postinsts.list
/etc/init.d/run-postinsts
/usr/share/run-postinsts/run-postinsts.awk
/usr/share/run-postinsts/
