# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/state-ctrl.postinst
[ -n "$D" ] && OPT="-r $D" || OPT="-s"

    # generate init scripts
    update-rc.d $OPT -f check_device_code remove
    update-rc.d $OPT check_device_code start 40 S .
    update-rc.d $OPT -f start_state_mgr remove
    update-rc.d $OPT start_state_mgr defaults 10
    update-rc.d $OPT -f start_off_chg remove
    update-rc.d $OPT start_off_chg defaults 19
    update-rc.d $OPT swap_on defaults 99
