# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/storageshare.list
/etc/default_config/storageshare
/usr/bin/change_storageshare_config
/etc/init.d/service_storageshare
/etc/config/storageshare
