# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/system-core-liblog.postinst
[ -n "$D" ] && OPT="-r $D" || OPT="-s"
        update-rc.d $OPT -f adbd remove
        update-rc.d $OPT adbd start 42 S 2 3 4 5 S . stop 80 0 1 6 .
        update-rc.d $OPT -f usb remove
        update-rc.d $OPT usb start 41 S .
