# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/tpoled-s90319.list
/etc/init.d/start_oledd
/etc/oled_animation.config
/etc/oled_res
/usr/bin/oledd
