# ghidra_projects/format_string/M7350(EU) 3.0_1.1.1 Build 160330 Rel.1002n_User.zip/system/var/lib/opkg/info/update.list
/etc/default_config/SystemUpdate
/usr/bin/prepareUdisk
/usr/bin/updatecheck
/usr/fota/
/etc/init.d/upgraded
/etc/config/SystemUpdate
/usr/lib/libupgrade.so
